# resources/supervisor-minion.conf
[program:salt-minion]
command=/usr/bin/salt-minion
autostart=true
autorestart=true
stderr_logfile=/var/log/salt-minion.log
stdout_logfile=/var/log/salt-minion.err
stopasgroup=true